# systemd/linein-passthrough-refresh.service
[Unit]
Description=Refresh line-in PipeWire loopback config
After=pipewire.service wireplumber.service
Wants=pipewire.service wireplumber.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/linein-passthrough-refresh

[Install]
WantedBy=default.target
